# config/cah/aragorn/etc/NetworkManager/dispatcher.d/apply-cah-aragorn-tc-rules
../../../../../../nm-dispatchers/apply-tc-rules