# config/home/aldarion/etc/NetworkManager/dispatcher.d/nautilus-de-vpn-autoconnect
../../../../../../nm-dispatchers/nautilus-de-vpn-autoconnect