# config/home/elrond/etc/NetworkManager/dispatcher.d/cah-vpn-autoconnect
../../../../../../nm-dispatchers/cah-vpn-autoconnect