# config/home/elrond/etc/NetworkManager/dispatcher.d/ruckus-vpn-autoconnect
../../../../../../nm-dispatchers/ruckus-vpn-autoconnect